# packages/tp-docker/DEBIAN/postinst
#!/bin/bash
set -e

# Get the user who invoked sudo (or current user if not using sudo)
ACTUAL_USER="${SUDO_USER:-${USER}}"

# If we're root and no SUDO_USER is set, try to find the default user
if [ "$ACTUAL_USER" = "root" ] || [ -z "$ACTUAL_USER" ]; then
    # Look for user with UID 1000 (default WSL user)
    ACTUAL_USER=$(getent passwd 1000 | cut -d: -f1)
fi

# Ensure Docker CE repository and packages are installed
if ! command -v docker > /dev/null 2>&1; then
    echo "Configuring official Docker CE repository..."
    install -m 0755 -d /etc/apt/keyrings 2>/dev/null || true
    curl -fsSL --connect-timeout 10 --max-time 30 https://download.docker.com/linux/ubuntu/gpg -o /etc/apt/keyrings/docker.asc 2>/dev/null || true
    chmod a+r /etc/apt/keyrings/docker.asc 2>/dev/null || true

    # Fallback resolute to noble since download.docker.com publishes LTS releases
    CODENAME=$(. /etc/os-release 2>/dev/null && echo "$VERSION_CODENAME")
    if [ "$CODENAME" = "resolute" ] || [ -z "$CODENAME" ]; then
        CODENAME="noble"
    fi

    if [ -f /etc/apt/keyrings/docker.asc ]; then
        echo "deb [arch=$(dpkg --print-architecture) signed-by=/etc/apt/keyrings/docker.asc] https://download.docker.com/linux/ubuntu $CODENAME stable" > /etc/apt/sources.list.d/docker.list
        apt-get update >/dev/null 2>&1 || true
    fi

    echo "Installing Docker CE packages..."
    DEBIAN_FRONTEND=noninteractive apt-get install -y docker-ce docker-ce-cli containerd.io docker-buildx-plugin docker-compose-plugin 2>/dev/null || DEBIAN_FRONTEND=noninteractive apt-get install -y docker.io 2>/dev/null || true
fi

# Create docker group if it doesn't exist
if ! getent group docker > /dev/null 2>&1; then
    echo "Creating docker group..."
    groupadd docker
fi

# Add user to docker group and configure ACLs for immediate sudo-less execution
if [ -n "$ACTUAL_USER" ] && [ "$ACTUAL_USER" != "root" ]; then
    echo "Adding user '$ACTUAL_USER' to docker group..."
    usermod -aG docker "$ACTUAL_USER" 2>/dev/null || true

    # Grant POSIX ACL on /var/run/docker.sock if setfacl is available
    if [ -e /var/run/docker.sock ] && command -v setfacl >/dev/null 2>&1; then
        echo "Setting socket ACL for immediate execution..."
        setfacl -m u:"$ACTUAL_USER":rw /var/run/docker.sock 2>/dev/null || true
    fi

    # Configure persistent systemd tmpfiles.d rule to maintain ACL on socket creation
    mkdir -p /etc/tmpfiles.d
    cat > /etc/tmpfiles.d/docker-acl.conf << EOF
# Grant POSIX ACL read/write permissions on docker socket to $ACTUAL_USER
a+ /var/run/docker.sock - - - - u:$ACTUAL_USER:rw
EOF
else
    echo "Warning: Could not determine user to add to docker group."
    echo "Please run manually: sudo usermod -aG docker \$USER"
fi

# Enable and start Docker service if systemd is active, or configure for next boot
if [ -d /run/systemd/system ]; then
    echo "Enabling and starting Docker service via systemd..."
    systemctl enable docker 2>/dev/null || true
    systemctl start docker 2>/dev/null || echo "Note: Docker service will start on next boot"

    # Re-apply setfacl after service start if socket exists
    if [ -n "$ACTUAL_USER" ] && [ "$ACTUAL_USER" != "root" ] && [ -e /var/run/docker.sock ] && command -v setfacl >/dev/null 2>&1; then
        setfacl -m u:"$ACTUAL_USER":rw /var/run/docker.sock 2>/dev/null || true
    fi
else
    echo "System is not currently running systemd (PID 1)."
    echo "Enabling Docker service for systemd boot..."
    mkdir -p /etc/systemd/system/multi-user.target.wants
    ln -sf /lib/systemd/system/docker.service /etc/systemd/system/multi-user.target.wants/docker.service 2>/dev/null || true
    update-rc.d docker enable 2>/dev/null || true
fi

echo ""
echo "Docker CE installation complete!"
echo ""
echo "Available commands:"
echo "  docker --version"
echo "  docker compose version"
echo "  docker buildx version"
echo ""

exit 0
